# post-metrics.service
[Unit]
Description=Post system metrics

[Service]
ExecStart=/bin/post-metrics
