# post-metrics.timer
[Unit]
Description=Post metrics every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
